# NMDesktop/.config/eww/eww.service
[Unit]
Description=eww widgets

[Service]
ExecStart=eww daemon && eww open nmdesktop
ExecStop=killall eww